# dynamic-layers/flutter-layer/recipes-unit/flutter-sample-no-material/systemd-fsnm/fsnm.service
[Unit]
Description=Boot to flutter sample no material

[Service]
ExecStart=/usr/bin/flutter-pi --release /usr/share/flutter/flutter-sample-no-material/@@FLUTTER_SDK_TAG@@/release

[Install]
WantedBy=multi-user.target
